# printer-proxy/printer-proxy.service
[Unit]
Description=Printer Proxy Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/printer-proxy
ExecStart=/usr/bin/node /root/printer-proxy/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
